# Biblioteca/config/supervisor/conf.d/supervisor.conf
[program:Bilioteca]
command=/config/gunicorn Bilioteca.wsgi:application --bind 0.0.0.0:8000
directory=/app
autostart=true
autorestart=true
redirect_stderr=true